# Library/LaunchAgents/backup.rclone.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>backup.rclone</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/fish</string>
        <string>-c</string>
        <string>~/cron/rclone_backup.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>86400</integer>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/backup.rclone.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/backup.rclone.stderr.log</string>
</dict>
</plist>
